# etc/crontab
# Add to /etc/crontab or another crontab to bootstrap the periodic job run.

# Perform maintenance every minute.
*	*	*	*	*	root	/opt/bstg/bin/periodicf minute
